# build/supervisor/conf.d/postfix.ini
[program:postfix]
command=/usr/local/bin/postfix.sh
