# build/supervisor/conf.d/sympa-fcgi.ini
[program:sympa-fcgi]
command=/usr/lib/cgi-bin/sympa/wwsympa.fcgi
environment=FCGI_SOCKET_PATH="127.0.0.1:9000"

user=sympa
autostart=true
autorestart=true