# build/supervisor/conf.d/sympa-bounced.ini
[program:sympa-bounced]
command=/usr/lib/sympa/bin/task_manager.pl --foreground

user=sympa
autostart=true
autorestart=true

; [Unit]
; Description=SYMPA mailing list manager (task_manager)
; PartOf=sympa.service
; Documentation=man:task_manager(8)

; [Service]
; User=sympa
; Group=sympa

; PermissionsStartOnly=true
; ExecStartPre=-/bin/mkdir -p /run/sympa
; ExecStartPre=-/bin/chown sympa:sympa /run/sympa

; Type=forking
; ExecStart=/usr/lib/sympa/bin/task_manager.pl
; PIDFile=/run/sympa/task_manager.pid

; Restart=always
; RestartSec=15

; PrivateTmp=true

; [Install]
; WantedBy=multi-user.target
; root@0c230dd8ab75:/# cat /lib/systemd/system/sympa                      
; sympa-archived.service      sympa-bounced.service       sympa-bulk.service          sympa-task_manager.service  sympa.service               
; root@0c230dd8ab75:/# cat /lib/systemd/system/sympa-bounced.service 
; [Unit]
; Description=SYMPA mailing list manager (bounced)
; PartOf=sympa.service
; Documentation=man:bounced(8)

; [Service]
; User=sympa
; Group=sympa

; PermissionsStartOnly=true
; ExecStartPre=-/bin/mkdir -p /run/sympa
; ExecStartPre=-/bin/chown sympa:sympa /run/sympa

; Type=forking
; ExecStart=/usr/lib/sympa/bin/bounced.pl
; PIDFile=/run/sympa/bounced.pid

; Restart=always
; RestartSec=15

; PrivateTmp=true

; [Install]
; WantedBy=multi-user.target