# build/supervisor/conf.d/sympa-task_manager.ini
[program:sympa-task_manager]
command=/usr/lib/sympa/bin/task_manager.pl --foreground

user=sympa
autostart=true
autorestart=true

; [Unit]
; Description=SYMPA mailing list manager (task_manager)
; PartOf=sympa.service
; Documentation=man:task_manager(8)

; [Service]
; User=sympa
; Group=sympa

; PermissionsStartOnly=true
; ExecStartPre=-/bin/mkdir -p /run/sympa
; ExecStartPre=-/bin/chown sympa:sympa /run/sympa

; Type=forking
; ExecStart=/usr/lib/sympa/bin/task_manager.pl
; PIDFile=/run/sympa/task_manager.pid

; Restart=always
; RestartSec=15

; PrivateTmp=true

; [Install]
; WantedBy=multi-user.target